# backend/crontab
* */3 * * * python3 /root/reviewgram/cron/clear_tokens.py
