# backend/crontab
* */3 * * * python3 /root/reviewgram/cron/clear_tokens.py
*/6 * * * * python3 /root/reviewgram/cron/clear_chats.py
